# deploy/sub2freeApi.service
[Unit]
Description=Sub2FreeApi - AI API Gateway Platform
Documentation=https://github.com/Wei-Shaw/sub2api
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=sub2freeapi
Group=sub2freeapi
WorkingDirectory=/home/third_party/sub2freeApi/deploy
EnvironmentFile=/home/third_party/sub2freeApi/deploy/sub2api.env
ExecStart=/home/third_party/bin/sub2freeApi/sub2freeApi
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sub2freeApi

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
PrivateTmp=true
ReadWritePaths=/home/third_party/sub2freeApi/deploy/data

# Environment - Server configuration
# Modify these values to change listen address and port
Environment=GIN_MODE=release
Environment=DEPLOYMENT_PROFILE=free
Environment=SERVER_HOST=0.0.0.0
Environment=SERVER_PORT=18382

[Install]
WantedBy=multi-user.target
